# deploy/ec2/color-game-server.service.example
[Unit]
Description=Color Game Socket Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/srv/color-game/COLOR-GAME
EnvironmentFile=/etc/color-game/server.env
ExecStart=/usr/bin/pnpm --filter @color-game/server start
Restart=always
RestartSec=5
KillSignal=SIGTERM
SyslogIdentifier=color-game-server

[Install]
WantedBy=multi-user.target
